# root/etc/s6-overlay/s6-rc.d/init-foldingathome-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions on config root and folders
lsiown -R abc:abc \
    /config
